# supervisord.conf

[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20

[program:init-wp-content]
command=/bin/sh -c "if [ ! -d /var/www/html/wp-content/plugins ]; then cp -r /docker-wp-content/* /var/www/html/wp-content/; fi"
priority=1
autostart=true
autorestart=false
